# debian/unity-settings-daemon.user-session.upstart
description "The Unity Settings Daemon"
author "Stéphane Graber <stgraber@ubuntu.com>"

start on started dbus and starting gnome-session and xsession SESSION=ubuntu
stop on stopping gnome-session DESKTOP_SESSION=ubuntu

respawn
exec /usr/lib/unity-settings-daemon/unity-settings-daemon
